# deploy/hairyclaw.service
[Unit]
Description=HairyClaw Agent (%i)
After=network-online.target ollama.service
Wants=network-online.target

[Service]
Type=simple
# WorkingDirectory and EnvironmentFile should be set per-deployment in local.toml or override
WorkingDirectory=%h/hairyclaw
EnvironmentFile=%h/hairyclaw/.env
ExecStart=/usr/bin/node apps/hairy-agent/dist/main.js
Restart=on-failure
RestartSec=10

# Resource limits
MemoryMax=2G
CPUQuota=200%

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=hairyclaw

[Install]
WantedBy=default.target
